# 01-detection-engineering/p6-threat-informed/rules/win_scheduled_task_creation.yml
title: Scheduled Task Creation via schtasks
id: 8a1b2c3d-0004-4a00-9000-000000000004
status: experimental
description: Detects scheduled task creation using schtasks, commonly abused for persistence.
author: portfolio
references:
  - https://attack.mitre.org/techniques/T1053/005/
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
    Image|endswith: '\schtasks.exe'
  selection_cli:
    CommandLine|contains: '/create'
  condition: selection_img and selection_cli
falsepositives:
  - Administrators and software create scheduled tasks routinely.
level: low
tags:
  - attack.persistence
  - attack.t1053.005
